# services/core/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/logs/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:wireguard]
command=/scripts/core/start-wireguard.sh
autostart=true
autorestart=true
startretries=3
stdout_logfile=/logs/wireguard.log
stderr_logfile=/logs/wireguard.error.log
priority=10

[program:sshd]
command=/usr/sbin/sshd -D -e
autostart=true
autorestart=true
startretries=3
stdout_logfile=/logs/sshd.log
stderr_logfile=/logs/sshd.error.log
priority=20
depends_on=wireguard

[program:restic]
command=/scripts/core/start-restic.sh
autostart=true
autorestart=true
startretries=3
stdout_logfile=/logs/restic.log
stderr_logfile=/logs/restic.error.log
priority=30
depends_on=wireguard,sshd
